# webapp/scripts/itrequest-frontend.service
[Unit]
Description=IT Request Platform — Frontend (vite preview)
After=network.target itrequest-backend.service

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__APP_DIR__/webapp/frontend
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 5173
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
